# ops/systemd/saytu-yef-backup.timer
[Unit]
Description=Sauvegarde quotidienne de Saytu Yef

[Timer]
OnCalendar=*-*-* 02:15:00 UTC
Persistent=true
RandomizedDelaySec=15m

[Install]
WantedBy=timers.target
